# debian/herald.postinst
#!/bin/bash

# herald.postinst
#
# Copyright (C) 2019 Kano Computing Ltd.
# License: http://www.gnu.org/licenses/gpl-2.0.txt GNU GPLv2
#
# Post debian package install script.


case "$1" in
    configure)
        # Enable system systemd services.
        systemctl enable start-user-internet-available.service
        ;;
esac
